# lib/slacksend.cron
*/2 *     * * *     www-data   [ -x /var/www/cyphre/apps/slacknotify/lib/cron.php ] && /var/www/cyphre/apps/slacknotify/lib/cron.php
